# launch/com.gideon.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gideon.web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/rickberghuis/Documents/Claude assistent/.venv/bin/python</string>
    <string>-m</string>
    <string>gideon.main</string>
    <string>--web</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/rickberghuis/Documents/Claude assistent</string>

  <!-- Start at login, and restart if it ever exits. -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Logs land in the git-ignored state/ dir. -->
  <key>StandardOutPath</key>
  <string>/Users/rickberghuis/Documents/Claude assistent/state/web.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/rickberghuis/Documents/Claude assistent/state/web.log</string>
</dict>
</plist>
